# admin/supervisord.conf
[supervisord]
nodaemon=true

[program:admin-app]
command=python3 manage.py runserver 0.0.0.0:8000
directory=/app
autostart=true
autorestart=true

[program:consumer]
command=python3 consumer.py
directory=/app
autostart=true
autorestart=true
